# kairos/overlays/role/storage/oci/etc/systemd/system/rtslib-fb-targetctl.service.d/10-k2-after-zfs.conf
# targetctl restore fails for block backstores whose backing device does
# not exist yet, so the LIO config restore must wait for the ZFS pool
# import and for /dev/zvol/* links (zvol_wait). Wants= rather than
# Requires=: a freshly installed appliance has no pool yet, and restoring
# an empty LIO config must still succeed.
[Unit]
After=zfs.target zfs-volume-wait.service
Wants=zfs-volume-wait.service
