# kairos/overlays/role/storage/oci/etc/systemd/system/rtslib-fb-targetctl.service.d/10-k2-after-zfs.conf
# targetctl restore fails for block backstores whose backing device does
# not exist yet, so the LIO config restore must wait for the ZFS pool
# import, key load, and /dev/zvol/* links (zvol_wait). Wants= rather
# than Requires=: a freshly installed appliance has no pool yet, and
# restoring an empty LIO config must still succeed. The chain is
# import -> k2-zfs-load-key -> zfs-volume-wait -> LIO restore.
[Unit]
After=zfs.target zfs-volume-wait.service
Wants=zfs-volume-wait.service
