# kairos/overlays/role/storage/oci/etc/systemd/system/k2-storage-health.service
[Unit]
Description=K2 storage appliance health check
After=zfs.target rtslib-fb-targetctl.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/k2-node-agent storage-health
