# kairos/overlays/role/storage/oci/etc/systemd/system/k2-storage-health.timer
[Unit]
Description=Periodic K2 storage appliance health check

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
