# kairos/image-build/overlays/kubernetes/k3s/oci/etc/systemd/system/k3s.service.d/wait-for-time-sync.conf
# K2: hard-require system time has been synchronized with NTP before
# k3s starts. The CM4 has no on-board RTC; on cold boot the kernel
# image's build date becomes the initial wall clock. Without this gate,
# k3s + etcd record Node creationTimestamps in the wrong era AND, more
# dangerously, every x509 cert minted while time is off (kubelet,
# apiserver, etcd peer, kube-vip) carries a notBefore in the wrong era
# — TLS handshakes silently fail later when the real clock advances
# past the cert's "valid from" window in unexpected ways.
#
# Requires= (hard) rather than Wants=: if NTP is unreachable, we'd
# rather k3s refuses to start than come up with a clock that will
# corrupt TLS material across the cluster. The OnFailure handler on
# systemd-time-wait-sync (sibling drop-in) emits a K2-tagged warning to
# journal so the operator sees the situation immediately.
[Unit]
After=time-sync.target
Requires=time-sync.target
