# kairos/overlays/role/k8s/oci/etc/systemd/system/systemd-time-wait-sync.service.d/timeout.conf
# K2: cap the NTP-sync wait at 60s and route any failure to a clear
# K2-tagged log entry so operators can distinguish "this host can't
# reach NTP" from generic boot delays.
#
# 60s is a balance — long enough for a sane DHCP+DNS+NTP handshake on
# a healthy network, short enough that a misconfigured firewall doesn't
# stretch boot beyond the operator's patience. systemd marks the unit
# 'failed' on timeout; the OnFailure handler then logs a K2-tagged
# warning to journal (grep with `journalctl -t k2`).
[Service]
TimeoutStartSec=60s

[Unit]
OnFailure=k2-time-sync-warning.service
