# kairos/overlays/role/k8s/oci/etc/systemd/system/k3s-agent.service.d/wait-for-time-sync.conf
# K2: hard-require NTP sync before k3s-agent (worker role) starts.
# Same x509-cert rationale as the sibling drop-in for k3s.service —
# kubelet TLS bootstrap certs minted while the clock is wrong will
# silently break later. Requires= (hard) rather than Wants=.
[Unit]
After=time-sync.target
Requires=time-sync.target
